# supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/log/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/tmp/supervisord.pid

[program:mongodb]
command=mongod --dbpath /data/db --logpath /data/log/mongod.log --bind_ip 127.0.0.1 --port 27017
autostart=true
autorestart=true
priority=10
startsecs=5
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:api]
command=/app/api
directory=/app
autostart=true
autorestart=true
priority=20
startsecs=3
startretries=15
environment=
    MONGO_URI="mongodb://%(ENV_MONGO_APP_USER)s:%(ENV_MONGO_APP_PASS)s@127.0.0.1:27017/vulnsamurai",
    MONGO_DB="vulnsamurai",
    JWT_SECRET="%(ENV_JWT_SECRET)s",
    BIND_ADDR="127.0.0.1:8000",
    FRONTEND_ORIGIN="http://localhost:3000",
    RUST_LOG="vuln_samurai_api=info,tower_http=info"
stdout_logfile=/data/log/backend.log
stderr_logfile=/data/log/backend.log

[program:frontend]
command=node server.js
directory=/app/frontend
autostart=true
autorestart=true
priority=30
startsecs=3
startretries=10
stdout_logfile=/data/log/frontend.log
stderr_logfile=/data/log/frontend.log
